# systemd/agentic-cloud-sleep-guard.service
[Unit]
Description=Detach agentic-cloud rclone FUSE mounts around system sleep
Documentation=https://github.com/Asi0Flammeus/agentic-cloud
Before=sleep.target
StopWhenUnneeded=yes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/agentic-cloud-sleep-guard pre
ExecStop=/usr/local/sbin/agentic-cloud-sleep-guard post

[Install]
WantedBy=sleep.target
